# daemon/claude-strategist-market.timer
[Unit]
Description=Claude Strategist — Market Analysis Timer (every 2 hours)

[Timer]
OnCalendar=*-*-* 0/2:00:00
Persistent=true

[Install]
WantedBy=timers.target
