# daemon/claude-strategist-review.service
[Unit]
Description=Claude Strategist — Trade Review Cron
After=network.target

[Service]
Type=oneshot
User=strategist
WorkingDirectory=/opt/ArkNode-AI/projects/trading-signal-ai/strategist
ExecStart=/home/strategist/.bun/bin/bun run src/cron/trade-review.ts

# Environment
Environment=HOME=/home/strategist
Environment=PATH=/home/strategist/.bun/bin:/home/strategist/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/opt/ArkNode-AI/projects/trading-signal-ai/strategist/.env

# Timeout (Claude CLI can take up to 3 minutes)
TimeoutStartSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-strategist-review
